# etc/sv/mailcatcher/run
#!/bin/sh

exec 2>&1

#exec mailcatcher --smtp-ip 127.0.0.1 --http-ip 0.0.0.0 -f
exec sudo -i -u pair mailcatcher --smtp-ip 127.0.0.1 --http-ip 0.0.0.0 -f
